# focas/bin/arm/init.rc
#!/bin/sh /etc/rc.common
# Copyright (C) 2018 freeioe.org

START=90

USE_PROCD=1

FOCAS_UBUS_ROOT=/usr/focas_armhf_rootfs/sysroot
FOCAS_UBUS_BIN=focas_ubus

start_service () {
	procd_open_instance
	procd_set_param command $FOCAS_UBUS_ROOT/bin/arch-chroot "$FOCAS_UBUS_ROOT" "$FOCUS_UBUS_BIN"
	procd_set_param respawn
	procd_close_instance
}

